# systemd/bascula-ui.service.d/10-tty.conf
[Service]
StandardInput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
