# systemd/bascula-ap-ensure.service
[Unit]
Description=Ensure Bascula AP when no saved Wi-Fi profiles
After=NetworkManager.service

[Service]
Type=oneshot
TimeoutStartSec=20
ExecStartPre=-/usr/bin/nm-online -s -q -t 10
ExecStart=/opt/bascula/current/scripts/bascula-ap-ensure.sh

[Install]
WantedBy=multi-user.target
