# systemd/bascula-ap-ensure.service
[Unit]
Description=Ensure Bascula AP when no other connectivity
After=NetworkManager.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/bascula/current/scripts/bascula-ap-ensure.sh
User=root
Group=root
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
