# systemd/bascula-ocr.service
[Unit]
Description=Báscula OCR Service
After=network.target bascula-backend.service

[Service]
Type=simple
WorkingDirectory=/opt/bascula/current
ExecStart=/opt/bascula/current/.venv/bin/python -m backend.ocr_service
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
