# systemd/bascula-ap-ensure.service
[Unit]
Description=Ensure Bascula AP when no other connectivity
After=NetworkManager.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bascula-ap-ensure.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
