# systemd/bascula-backend.service
[Unit]
Description=Bascula Backend - FastAPI Server
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bascula/current
Environment="PATH=/opt/bascula/current/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/bascula/current/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8081 --workers 2 --log-level info
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
