# systemd/bascula-backend.service
[Unit]
Description=Bascula Backend (core)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bascula/current
ExecStart=/opt/bascula/current/.venv/bin/python -m backend.main --port 8081
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
