# systemd/bascula-health-wait.service
[Unit]
Description=Espera a que miniweb (8080) y backend (8081) estén disponibles
After=bascula-miniweb.service bascula-backend.service
Wants=bascula-miniweb.service bascula-backend.service

[Service]
Type=oneshot
ExecStart=/opt/bascula/current/systemd/bin/bascula-health-wait.sh

[Install]
WantedBy=multi-user.target
