# systemd/bascula-miniweb.service
[Unit]
Description=Bascula Mini-Web Backend
After=NetworkManager.service NetworkManager-wait-online.service bascula-ap-ensure.service
Wants=NetworkManager-wait-online.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bascula/current
Environment=PATH=/opt/bascula/current/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=BASCULA_CFG_DIR=/home/pi/.bascula
Environment=BASCULA_CAM_RELEASE_ON_EXIT=true
ExecStartPre=-/usr/bin/nm-online -s -q -t 25
ExecStart=/opt/bascula/current/.venv/bin/python -m uvicorn backend.miniweb:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
