# systemd/bascula-backend.service
[Unit]
Description=Bascula Backend (core)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bascula/current
Environment=BASCULA_BACKEND_HOST=0.0.0.0
Environment=BASCULA_BACKEND_PORT=8081
ExecStart=/opt/bascula/current/.venv/bin/python -m backend.main --host ${BASCULA_BACKEND_HOST} --port ${BASCULA_BACKEND_PORT}
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
