# systemd/bascula-ap-ensure.service
[Unit]
Description=Ensure Bascula AP when no other connectivity
After=NetworkManager.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# No bloquear: si nm-online falla, continuar igualmente
ExecStartPre=-/usr/bin/nm-online -s -q -t 25
ExecStart=/opt/bascula/current/scripts/bascula-ap-ensure.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
