# systemd/bascula-ui.service
[Unit]
Description=Bascula Digital Pro UI (Kiosk)
After=graphical.target network.target

[Service]
User=pi
Type=simple
Environment=DISPLAY=:0
ExecStart=/opt/bascula/current/scripts/start-kiosk.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
