# systemd/bascula-ap-ensure.service
[Unit]
Description=Ensure Bascula AP when no other connectivity
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/opt/bascula/current/scripts/bascula-ap-ensure.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
