# packaging/systemd/bascula-miniweb.service
[Unit]
Description=Bascula Mini-Web Backend
After=network.target NetworkManager.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bascula/current
Environment=PATH=/opt/bascula/current/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=BASCULA_CFG_DIR=/home/pi/.bascula
# Permite consultar el PIN localmente y durante modo AP (la API decide). 0 para desactivar.
Environment=BASCULA_ALLOW_PIN_READ=1
ExecStart=/opt/bascula/current/.venv/bin/python -m uvicorn backend.miniweb:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
